# systemd/dci-automatic-deploy.timer
[Unit]
Description=dci-auto-deploy timer

[Timer]
OnCalendar=13:37
Persistent=true     

[Install]
WantedBy=basic.target
